# xe-devel/opt/xundeenergie/etc/systemd/system/apt-daily.service.d/service.conf
[Service]
ExecStartPre=/usr/bin/nm-updatecheck.py -G -c

[Install]
WantedBy=timer-daily.target

